# python_app/crontab
0 0 * * * root /app/etl.sh >> /var/log/cron.log 2>&1
